# services/typo3-v7-tech-autodeploy.path
[Unit]
Description=Path watcher for autodeploy of TYPO3 (v7) artifacts for development branches

[Path]
DirectoryNotEmpty=/srv/artifacts/typo3-v7/tech

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=role=fdc