# services/typo3-live-autodeploy.path
[Unit]
Description=Path watcher for autodeploy of TYPO3 artifacts for live environment

[Path]
DirectoryNotEmpty=/srv/artifacts/typo3/live

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=role=bil